# tools/print-agent/build/macos/com.ironwaves.print-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ironwaves.print-agent</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/ironwaves-print-agent/ironwaves-print-agent</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/com.ironwaves.print-agent.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/com.ironwaves.print-agent.err.log</string>
</dict>
</plist>
